# projects/pocketchip/openwrt/files/etc/init.d/usb-gadget
#!/bin/sh /etc/rc.common
# USB gadget bring-up. The actual work is done by the kernel via
# /etc/modules.d/50-g_cdc, which autoloads the g_cdc legacy gadget
# (CDC-ACM + CDC-ECM, both standards macOS / Linux / Windows recognise
# natively) at boot. This script just logs the resulting state so
# 'logread -t usb-gadget' tells you whether the gadget bound to the UDC.
#
# If we later need a more elaborate composite (mass storage, RNDIS for
# Windows-only hosts, multiple ACM endpoints) we'll switch to a
# configfs-driven setup here. g_cdc is enough for the SSH/console use
# case this project actually has.

START=11
USE_PROCD=1

start_service() {
	if [ -d /sys/class/udc ] && [ -n "$(ls /sys/class/udc 2>/dev/null)" ]; then
		UDC=$(ls /sys/class/udc | head -1)
	else
		UDC="<none>"
	fi

	if [ -d /sys/class/net/usb0 ]; then
		logger -t usb-gadget "g_cdc bound to UDC $UDC; usb0 present"
	else
		logger -t usb-gadget "g_cdc did NOT bind (UDC=$UDC, usb0 missing)"
	fi
}
